# env-scripts/daily.cron
# m h  dom mon dow   command
# At 1 am run the task once every day.
00 01 * * * /bin/bash /usr/local/exec_daily.sh >> /usr/local/data/exec_daily.log 2>&1
